# device/pi/systemd/user/pisignage-vlc.service
[Unit]
Description=PiSignage Native VLC Video Player

[Service]
Type=simple
WorkingDirectory=/home/donnoel/PiSignage
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=WAYLAND_DISPLAY=wayland-0
Environment=PISIGNAGE_REPO_ROOT=/home/donnoel/PiSignage
Environment=PISIGNAGE_DISPLAY_OUTPUT=HDMI-A-1
Environment=PISIGNAGE_DISPLAY_RESOLUTION=1920x1080@60.000000
Environment=PISIGNAGE_VLC_VIDEO_OUTPUT=wl_shm
Environment=PISIGNAGE_VLC_PLAYBACK_MODE=continuous
Environment=PISIGNAGE_PLAYLIST_HANDOFF_OVERLAP_MS=2500
Environment=PISIGNAGE_VLC_RESTART_BACKOFF_MS=15000
Environment=PISIGNAGE_VLC_RESTART_BACKOFF_MAX_MS=120000
Environment=PISIGNAGE_VLC_WAYLAND_DISPLAY=wayland-0
Environment=PISIGNAGE_STATUS_PATH=%h/.local/state/pisignage/player-status.json
Environment=PISIGNAGE_STARTUP_SETTLE_MS=0
ExecStartPre=%h/.local/bin/pisignage-hide-desktop.sh
ExecStart=/usr/bin/node %h/.local/bin/pisignage-vlc-playlist.mjs
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
