# device/pi/systemd/user/pisignage-device-agent.service
[Unit]
Description=PiSignage Device Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/donnoel/PiSignage
Environment=PISIGNAGE_PLAYLIST_PATH=/home/donnoel/PiSignage/sample-content/playlist.local.json
Environment=PISIGNAGE_CACHE_DIR=%h/.local/cache/pisignage/device-agent
Environment=PISIGNAGE_HEARTBEAT_PATH=%h/.local/state/pisignage/heartbeat.json
Environment=PISIGNAGE_HEARTBEAT_INTERVAL_SECONDS=30
Environment=PISIGNAGE_HEARTBEAT_JITTER_SECONDS=5
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=WAYLAND_DISPLAY=wayland-0
EnvironmentFile=-%h/.config/pisignage/device-agent.env
ExecStart=/usr/bin/node /home/donnoel/PiSignage/device-agent/dist/index.js --loop
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
